# backend/data/odds_archive/odds-backup.timer
[Unit]
Description=Odds Archive Database Backup Timer (Daily at 3 AM)
Requires=odds-backup.service

[Timer]
# Run daily at 3:00 AM
OnCalendar=daily
OnCalendar=*-*-* 03:00:00

# Run 5 minutes after boot if missed
OnBootSec=5min

[Install]
WantedBy=timers.target
