# backend/data/odds_archive/odds-collector.timer
[Unit]
Description=Odds Archive Collection Timer (Every 30 minutes)
Requires=odds-collector.service

[Timer]
# Run every 30 minutes
OnBootSec=5min
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target
